# utilities/ddcci.timer
[Unit]
Description=Delay loading ddcci module

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
