# utilities/ddcci.service
[Unit]
Description=Load ddcci module

[Service]
Type=oneshot 
ExecStart=/usr/bin/modprobe ddcci_backlight
