# utilities/lowbat.service
[Unit]
Description=Hibernate on low battery

[Service]
ExecStart=/usr/local/bin/lowbat.sh
