# keyszer/etc/systemd/system/keyszer.target
[Unit]
Description=Keyszer

[Install]
WantedBy=graphical.target
