# keyszer/etc/systemd/system/uinput.service
[Unit]
Description=uinput

[Service]
ExecStart=/usr/bin/chmod 0660 /dev/uinput
Type=oneshot

[Install]
WantedBy=multi-user.target
