# systemd/etc/systemd/system/disable-bd-prochot.timer
[Unit]
Description=Disable BD PROCHOT timer

[Timer]
OnCalendar=minutely
Unit=disable-bd-prochot.service

[Install]
WantedBy=multi-user.target
