# systemd/etc/systemd/system/disable-bd-prochot.service
[Unit]
Description=Disable BD PROCHOT

[Service]
Type=oneshot
ExecStart=/usr/bin/disable-bd-prochot.sh
StandardOutput=journal
