# systemd/etc/systemd/system/redmibook-lid-close.service
[Unit]
Description=RedmiBook Lid Close Service

[Service]
ExecStart=redmibook-lid-close.sh
Restart=always

[Install]
WantedBy=multi-user.target
